# etc/rc.d/znc
#!/bin/sh

daemon="/usr/local/bin/znc"
daemon_user="katt"

. /etc/rc.d/rc.subr

rc_reload=NO
rc_cmd $1
